# laravel/docker/supervisord-staging.conf
[supervisord]
nodaemon=true

[program:cloudwatch]
process_name=%(program_name)s
command=/opt/aws/amazon-cloudwatch-agent/bin/start-amazon-cloudwatch-agent
autostart=true
autorestart=true
